# redis-checker.service
[Unit]
Description=VPN Redis Checker + Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/vpn_checker
ExecStart=/opt/vpn_checker/redis-checker -recheck -workers 10 -timeout 15s -serve 127.0.0.1:8081
Restart=on-failure
RestartSec=5s

Environment=REDIS_URL=redis://default:PASSWORD@HOST:PORT

# Logs: journalctl -u redis-checker -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
